# crontab
#
*/5 * * * * /usr/local/bin/python /app/parser_worker.py
#
